# digits/etc/supervisord.d/jupyter.conf
[program:jupyter]
command=/opt/conda/bin/ipython notebook --ip='*'
directory=/root/
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/jupyter.log

